# screenlock.target
[Unit]
Description=Lock screen
Before=suspend.target
StopWhenUnneeded=true

[Install]
RequiredBy=suspend.target
